# cookbooks/app_beta/files/default/beta.conf
[program:beta]
directory = /apps/beta/current
command = bundle exec rackup -p 9293 -P ./tmp/pids/rack.pid
user = liftopian
autostart = true