# cookbooks/app_beta/files/default/beta.conf
[program:beta]
directory = /apps/beta/current
command = bundle exec unicorn_rails -c ./config/unicorn.rb -E development
user = liftopian
autostart = true