# cookbooks/app_alpha/files/default/alpha.conf
[program:alpha]
directory = /apps/alpha/current
command = bundle exec unicorn_rails -c ./config/unicorn.rb -E production
user = liftopian
autostart = true
