# cookbooks/app_alpha/files/default/alpha.conf
[program:alpha]
directory = /apps/alpha/current
command = bundle exec rackup -P ./tmp/pids/rack.pid
user = liftopian
autostart = true